# scripts/update.cron
15 4 * * * /var/www/court.nanddevelopment.com/scripts/update.sh
